# vapt-server.service
[Unit]
Description=VAPT Scanner Server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/vapt-scanner-project
EnvironmentFile=/opt/vapt-scanner-project/.env
ExecStart=/opt/vapt-scanner-project/venv/bin/uvicorn backend.app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
